# crontab
* * * * * /usr/local/bin/php /home/user/public_html/index.php worker/index
